# deploy/drones-rule-unlock.timer
[Unit]
Description=Clean journal entry locks

[Timer]
OnBootSec=10s
OnUnitActiveSec=2m

[Install]
WantedBy=timers.target