# deploy/dronesjournalttlrules.timer
[Unit]
Description=Execute journal rules every 2 minutes

[Timer]
OnBootSec=10s
OnUnitActiveSec=2m

[Install]
WantedBy=timers.target