# deploy/drones-clean-old-tokens.timer
[Unit]
Description=Clean old tokens from the database

[Timer]
OnBootSec=120s
OnUnitActiveSec=23h

[Install]
WantedBy=timers.target